# backup/crontab
0 0 * * * /usr/local/bin/backup.sh >> /var/log/cron.log 2>&1
# An empty line is required for the file to be valid